# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Ensure script is executable
        chmod +x /usr/local/sbin/setup-persistent-identity

        # Mask regenerate_ssh_host_keys
        if systemctl list-unit-files regenerate_ssh_host_keys.service >/dev/null 2>&1; then
            systemctl mask regenerate_ssh_host_keys.service >/dev/null 2>&1 || true
        fi

        # Reload systemd and enable service
        systemctl daemon-reload
        systemctl enable wlanpi-persistent-identity.service

        echo "wlanpi-persistent-identity installed. Persistence will activate on next boot."
        echo "To activate now: sudo systemctl start wlanpi-persistent-identity.service"
        ;;
esac

exit 0